# Misc/10-wardrive.rules.example
SUBSYSTEM=="ieee80211", ACTION=="add", ATTR{address}=="$MAC", ATTR{type}=="1", NAME="$INTERFACE_NAME", ENV{NM_UNMANAGED}="1", RUN+="systemctl start hcx@$INTERFACE_NAME"
